# services/ircbot/run
#!/bin/bash
set -euxo pipefail
cd /opt/ircbot
exec /opt/ircbot/venv/bin/python -m ircbot.ircbot
